# xob-aud.service
[Unit]
Description=Brightness Bar

[Service]
Type=simple
ExecStart=/usr/bin/chainos-xob-aud
Restart=on-failure

[Install]
WantedBy=graphical-session.target
